# rootfs/etc/supervisor/conf.d/frpc.conf
[program:frpc]
command = /usr/bin/frpc -c /etc/frpc/frpc.ini
autorestart = true 
stderr_logfile = /var/log/frpc/stderr.log
stdout_logfile = /var/log/frpc/stdout.log
